# 99-mountcopy.rules
ACTION=="add", SUBSYSTEMS=="usb", KERNEL=="sd*", RUN+="/bin/systemctl start usb-mount@%k.service"
